# supervisor.ini
[program:capstone_gunicorn]
command=/srv/www/capstone/venv/bin/gunicorn -w 8 -b 127.0.0.1:8008 therefuge.wsgi:application
directory=/srv/www/capstone/django/
user=nobody
stdout_logfile=/srv/www/capstone/logs/gunicorn.log
stderr_logfile=/srv/www/capstone/logs/gunicornError.log
autostart=true
autorestart=true
priority=991
stopsignal=KILL